# udev/10-auto-ppp.rules

ACTION=="add", KERNEL=="ttyUSB*", ENV{DEVNAME}=="/dev/ttyUSB*", RUN+="/home/lfranchi/work/workspaces/autoroute/udev/udev-exec.sh"

